# cron.d/zabbix-kernelmisc
# [This file is part of the zabbix-agent-osso package]

# Check for "watchdog: BUG: soft lockup - CPU#24 stuck for 23s! [kvm:9643]"
# or "watchdog: watchdog0: watchdog did not stop!"
# or "mce: [Hardware Error]: CPU 15: Machine Check: 0 Bank 3: d820000000000150"
# or "AMD-Vi: Completion-Wait loop timed out"
#
# Note: when running this every 2 minutes, we see spurious errors like this:
#   An error was encountered while opening journal file or directory
#   /var/log/journal/ca7481b701f04ff390e4e896fccea062/system.journal,
#   ignoring file: No such file or directory
# We ignore those with the grep -v.
#
*/2 *    * * *  root  DST=/var/lib/zabbix/scripts/kernel.emerg && [ -x /usr/bin/journalctl ] && /usr/bin/journalctl -qo cat -S '-4 minutes' -p 1 _TRANSPORT=kernel 2>&1 | grep -vF ', ignoring file: No such file or directory' | tee $DST.log.tmp | wc -l >$DST.count.tmp && mv $DST.count.tmp $DST.count && sort $DST.log.tmp | uniq -c >$DST.log
