# cron.d/zabbix-smart
# [This file is part of the zabbix-agent-osso package]

# Update S.M.A.R.T. cache.
*/15 *   * * *  root  [ -x /etc/zabbix/scripts/smart.update_cache ] && /etc/zabbix/scripts/smart.update_cache
