# cron.d/zabbix-mysql2
# [This file is part of the zabbix-agent-osso package]

# We'll set the shell to bash and do a random sleep to avoid all machines
# running checks at the same exact time.
SHELL=/bin/bash

# Run a daily SELECT ... LIMIT 1; check on all databases.
# When the database consists of many smaller schemas, we'll run into issues
# with the duration of the script.
10 10 * * *  root  sleep $((RANDOM \% 1800)); for CLUSTER in $(/etc/zabbix/scripts/mysql2 list_instances); do DST=/var/lib/zabbix/scripts/mysql2.$CLUSTER; /etc/zabbix/scripts/mysql2-select-healthcheck $CLUSTER >$DST.tmp; mv $DST.tmp $DST; done
