# cron.d/zabbix-dpkg
# [This file is part of the zabbix-agent-osso package]

# Hourly update the outstanding security and regular updates.
35 */3     * * *  root  [ -x /etc/zabbix/scripts/dpkg.updates ] && /etc/zabbix/scripts/dpkg.updates
