# cron.d/zabbix-k8smgmt
# [This file is part of the zabbix-agent-osso package]

# Update K8S cached json files.
* *      * * *  root  [ -x /etc/zabbix/scripts/k8smgmt ] && /etc/zabbix/scripts/k8smgmt update-cache
