# cron.d/zabbix-dpkg
# [This file is part of the zabbix-agent-osso package]

# Hourly update the outstanding security and regular updates.
35 *     * * *  root  [ -x /usr/lib/update-notifier/update-motd-updates-available ] && /usr/lib/update-notifier/update-motd-updates-available --force >/dev/null
